# health-check.service
[Unit]
Description=Health Check Web Server
Documentation=https://github.com/your-org/health-check-server
After=network.target
Wants=network.target

[Service]
Type=simple
User=healthcheck
Group=healthcheck
WorkingDirectory=/etc/health-check-server
ExecStart=/usr/local/bin/health-check-server
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=health-check-server

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/etc/health-check-server /var/log/health-check-server
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true
RemoveIPC=true
PrivateDevices=true
ProtectKernelLogs=true
ProtectClock=true

# Network settings - allow outbound MQTT connections
# Note: Update IPAddressAllow if your MQTT broker is on a specific network
IPAddressDeny=any
IPAddressAllow=localhost/128 127.0.0.1/32
# Uncomment and add your MQTT broker IP if needed:
# IPAddressAllow=192.168.1.0/24

# Resource limits
LimitNOFILE=65536
MemoryMax=50M
CPUQuota=10%

[Install]
WantedBy=multi-user.target